# crontab
*/1 * * * * root python sim_worker.py
# An empty line is required at the end of this file for a valid cron file.
